# etc/runit/log/run
#!/bin/sh

LOGDIR=../logfiles

if ! test -d "$LOGDIR"; then
        mkdir -p "$LOGDIR"
fi

exec svlogd -tt "$LOGDIR"
